# files/nfdumpmetrics.cron
* * * * *   node-exp   /usr/local/sbin/nfdumpmetrics.sh
