# files/nflow-agg.cron
#!/bin/bash
#
# /etc/cron.daily/nflow-agg script
# daily aggregate and compress flows 
#

set -o allexport
source /etc/default/nfcapd
set +o allexport
[ -d "$NFLOW_DIR" ] || exit 1

cd $NFLOW_DIR
for file in `find . -name 'nfcapd.[0-9][0-9][0-9][0-9][0-9][0-9][0-9][0-9][0-9][0-9][0-9][0-9]'` ; do
    aggfile=$(echo ${file} | sed s/nfcapd/ajnfcapd/)
    nfdump -A srcip,dstip -j -a -r ${file} -w ${aggfile}
    [ -f ${aggfile} ] && rm -f ${file}
done
